# deploy/gh-desktop-receiver.service
[Unit]
Description=gh-desktop webhook receiver
Documentation=https://github.com/jpuckett11/gh-desktop
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ghd-recv
Group=ghd-recv
WorkingDirectory=/opt/gh-desktop-receiver
ExecStart=/opt/gh-desktop-receiver/.venv/bin/gh-desktop-receiver --host 127.0.0.1 --port 8765
Restart=on-failure
RestartSec=5
EnvironmentFile=/etc/gh-desktop-receiver/env
StateDirectory=gh-desktop-receiver
LogsDirectory=gh-desktop-receiver
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/gh-desktop-receiver /var/log/gh-desktop-receiver

[Install]
WantedBy=multi-user.target
